# xrdp/files/supervisord/xrdp-sesman.conf
[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman -n
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
